# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.4]
command=/usr/sbin/php-fpm7.4 -F

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
