# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm8.1]
command=/usr/sbin/php-fpm8.1 -F

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
